# django/supervisor.conf
[group:watchdog]
programs=goaccess

[program:goaccess]
command=/bin/bash django/goaccess.bash
directory=/project
autostart=yes
autorestart=yes
startretries=3
user=root
